# conf/supervisord-apache.conf
[program:apache2]
command=apache2 -D FOREGROUND
numprocs=1
autostart=true
autorestart=true